# web/systemd/viral.service.erb
[Unit]
Description=Viral
After=syslog.target
After=network.target

[Service]
ExecStart=/home/deployment/viral/start-viral.sh <%= @my_origin %>
KillMode=control-group

[Install]
WantedBy=multi-user.target
